# conf/systemd/xcal-device.path
[Unit]
Description=xcal-device monitor

[Path]
PathExists=/tmp/ip_acquired
Unit=xcal-device.service

[Install]
WantedBy=multi-user.target

